# config/systemd/user/spotifyd.service.d/override.conf
[Service]
Environment="PULSE_SINK=paspotify"