# config/systemd/user/irxevent.service
[Unit]
Description=irxevent

[Service]
Type=simple
ExecStart=/bin/irxevent

[Install]
WantedBy=default.target
