# config/systemd/user/socat@.service
[Install]
WantedBy=default.target

[Service]
Restart=on-success
ExecStart=/usr/bin/socat -u UDP-RECV:%i STDOUT
#514 usually

