# config/systemd/user/syncthing.service
[Unit]
Description=syncthing

[Service]
Type=simple
ExecStart=/bin/syncthing

[Install]
WantedBy=default.target
